# scrapy/cron/crontab
# Run Scrapy spider at 8:00 AM
0 8 * * * root /app/scrapy/scripts/scrape.sh >> /app/logs/cron.log 2>&1

# Ingest articles to Elasticsearch at 8:30 AM
30 8 * * * root /app/scrapy/scripts/ingest.sh >> /app/logs/cron.log 2>&1
